# systemd/emcp-reload.service
[Unit]
Description=Reload eMCP docker-compose stack
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/mnt/ops/docker/eMCP
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStartPost=/usr/bin/docker restart emcp-server
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
