# systemd/emcp-reload.path
[Unit]
Description=Watch for eMCP docker-compose reload trigger

[Path]
PathModified=/mnt/ops/docker/eMCP/.reload-trigger
Unit=emcp-reload.service

[Install]
WantedBy=multi-user.target
